# private_Library/private_LaunchAgents/com.chezmoi.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chezmoi.update</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>export PATH=/opt/homebrew/bin:/usr/local/bin:$PATH; chezmoi update --apply 2>&1 | logger -t chezmoi</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardErrorPath</key>
    <string>/tmp/chezmoi-update.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/chezmoi-update.log</string>
</dict>
</plist>
